# squashfs-root/etc/hotplug.d/iface/80-miqos-brguest
#!/bin/sh

[ "$ACTION" == "ifup" ] && [ "$INTERFACE" == "guest" ] && {
    /etc/init.d/miqos noifb_restart
}
